# files/ac.target
[Unit]
Description=On AC Power
DefaultDependencies=no
StopWhenUnneeded=yes

